# survey365/systemd/survey365.service
# Managed by survey365/install.sh — deployed to /etc/systemd/system/survey365.service
[Unit]
Description=Survey365 Field Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={user}
Group=dialout
WorkingDirectory={home}/rtk-surveying/survey365
ExecStart={home}/rtk-surveying/survey365/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3

# Environment
Environment=SURVEY365_DB={data_dir}/survey365.db
Environment=GNSS_PORT=/dev/ttyGNSS
Environment=GNSS_BAUD=115200
Environment=GNSS_BACKEND=ublox

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=survey365

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths={data_dir}
PrivateTmp=true

[Install]
WantedBy=multi-user.target
