# survey365/systemd/survey365-update.timer
[Unit]
Description=Survey365 periodic update check

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=survey365-update.service

[Install]
WantedBy=timers.target
