# survey365/systemd/survey365-boot.service
[Unit]
Description=Survey365 Boot Tasks (antenna voltage, auto-resume)
After=str2str_tcp.service
Wants=str2str_tcp.service

[Service]
Type=oneshot
ExecStart={home}/rtk-surveying/survey365/venv/bin/python3 {home}/rtk-surveying/survey365/app/boot.py
RemainAfterExit=yes

# Environment for boot script to find config and database
Environment=RTKBASE_DIR={home}/rtkbase
Environment=SURVEY365_DB={home}/rtk-surveying/survey365/data/survey365.db
Environment=SURVEY365_DIR={home}/rtk-surveying/survey365

# Runs as root for serial port access to F9P
User=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=survey365-boot

[Install]
WantedBy=multi-user.target
