# survey365/systemd/survey365-boot.service
# Managed by survey365/scripts/setup-pi.sh — deployed to /etc/systemd/system/survey365-boot.service
[Unit]
Description=Survey365 Boot Tasks (hardware check)
After=network.target

[Service]
Type=oneshot
ExecStart={home}/rtk-surveying/survey365/venv/bin/python3 {home}/rtk-surveying/survey365/app/boot.py
RemainAfterExit=yes

# Environment
Environment=SURVEY365_DB={data_dir}/survey365.db
Environment=GNSS_PORT=/dev/ttyGNSS

# Run as target user (dialout group provides serial access)
User={user}
Group=dialout

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=survey365-boot

[Install]
WantedBy=multi-user.target
