# systemd/survey365.service
# Managed by scripts/setup-pi.sh — deployed to /etc/systemd/system/survey365.service
[Unit]
Description=Survey365 Field Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={user}
Group=dialout
WorkingDirectory={repo_dir}
ExecStart={repo_dir}/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3
TimeoutStopSec=30

# Environment
Environment=SURVEY365_DB={data_dir}/survey365.db
Environment=GNSS_PORT=/dev/ttyGNSS
Environment=GNSS_BAUD=115200
Environment=GNSS_BACKEND=ublox

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=survey365

# Security hardening
# NoNewPrivileges is intentionally NOT set on the main service.
# survey365.service uses `sudo -n systemctl start/stop/reset-failed` to manage
# RTKLIB child units.  NoNewPrivileges=true prevents sudo from gaining the
# needed privileges.  This was verified on the Pi — enabling it breaks child
# unit control.  The RTKLIB child services *do* set NoNewPrivileges=true
# because they have no need for privilege escalation.
ProtectSystem=strict
ReadWritePaths={data_dir}
PrivateTmp=true

[Install]
WantedBy=multi-user.target
